# guest/inittab
::sysinit:/bin/mount -t proc  proc  /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/hostname scanner
::once:/usr/local/bin/run-scan.sh
::ctrlaltdel:/sbin/reboot